# docker/run/resources/supervisord.conf
[supervisord]
nodaemon=true

[program:xlr]
command=sh -x /opt/xlr/command.sh