# docker/run/resources/supervisord.conf
[supervisord]
nodaemon=true

[program:xld]
command=sh -x /opt/xlr/command.sh